# immo-panel.service
[Unit]
Description=ImmoMatch Control Panel
After=network.target

[Service]
WorkingDirectory=/app
ExecStart=/usr/bin/python3 /app/panel.py
Restart=always
RestartSec=3
Environment=PANEL_SECRET=immo2026

[Install]
WantedBy=multi-user.target
